# board/rpi/overlays/zero2w/etc/profile.d/xdg.sh
# mkdir -p /run/user/1000
# chmod 0700 /run/user/1000
# chown 1000:1000 /run/user/1000

export XDG_RUNTIME_DIR=/run/user/1000
